# server/root/etc/services.d/app/run
#!/usr/bin/with-contenv bash

cd /opt/app

echo "Starting data node server"
s6-setuidgid app python -m openapi_server